# work/systemd/notmuch.timer
[Unit]
Description=Notmuch Query Timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=default.target