# sudo/mnt-orange.automount
[Unit]
Description=Automount Orange Lacie Drive

[Automount]
Where=/mnt/orange

[Install]
WantedBy=multi-user.target
